# debian/cups-airprint.cron.d
@reboot		/usr/sbin/update-cups-airprint
